# system/stacks.target
[Unit]
Description=All Docker Compose Stacks
Requires=docker.service
After=docker.service network-online.target
